# ble_and_aws/pi-src/rooted-iot.service
[Unit]
Description=Rooted Robotics AWS IoT Connection
Documentation=https://github.com/rooted-robotics
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=rooted
Group=rooted
WorkingDirectory=/opt/rooted-ble
ExecStart=/opt/rooted-ble/.venv/bin/python3 -u /opt/rooted-ble/aws/command_handler.py
Restart=always
RestartSec=10

# Environment
Environment=PYTHONUNBUFFERED=1

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/opt/rooted-ble /home/rooted/te-cli

[Install]
WantedBy=multi-user.target
